# luci-app-netaudit/root/etc/init.d/netaudit
#!/bin/sh /etc/rc.common

USE_PROCD=1
PROG=/usr/lib/lua/netaudit/analyzer.lua
NAME=netaudit

START=99
STOP=10

boot() {
    return 0
}

start_service() {
    procd_open_instance
    procd_set_param command lua $PROG
    procd_set_param respawn
    procd_set_param respawn_retries 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_add_jail network
    procd_close_instance
    
    # 应用nftables规则
    if [ -f /etc/nftables.d/10-netaudit.nft ]; then
        nft -f /etc/nftables.d/10-netaudit.nft
    fi
    
    logger -t netaudit "Network Audit service started"
}

stop_service() {
    # 清理nftables规则
    nft delete table ip netaudit 2>/dev/null || true
    
    logger -t netaudit "Network Audit service stopped"
}

reload_service() {
    stop
    sleep 1
    start
}

restart() {
    stop
    sleep 2
    start
}
